# src/plm/plmcd/plmcd.service.in
[Unit]
Description=OpenSAF PLMC Daemon
After=local-fs.target remote-fs.target syslog.target network.target plmcboot.service
Before=opensafd.service
Wants=plmcboot.service

[Service]
ExecStart=@sbindir@/plmcd -c @sysconfdir@/plmcd.conf
Type=notify
Restart=always
NotifyAccess=main

[Install]
WantedBy=multi-user.target
