# osaf/services/infrastructure/fm/fms/scripts/fmd.in
#!/bin/sh
#
### BEGIN INIT INFO
# Provides: fmd 
# Required-Start: cron
# Should-Start: opensafd
# Required-Stop: cron
# Should-Stop: opensafd
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: 
# Description: 
### END INIT INFO

# Source LSB functions library
. /lib/lsb/init-functions
RETVAL=0

start() {
	return 0
}

stop() {
	logger -t fmd "Isolating the node."
	#tipc-config -lt=1.1.47:eth1-1.1.31:eth0/500
	
	pkill osafamfd
	pkill osafamfnd
	pkill osaffmd

	x=`tipc-config -b | grep eth`
	tipc-config -bd=$x
	
	return 0
}

case "$1" in
  start)
	start
	RETVAL=$?
	;;
  stop)
	stop
	RETVAL=$?
	;;
  *)
	RETVAL=2
esac

exit $RETVAL
