# osaf/services/infrastructure/fm/fms/scripts/osaf-fmd.in
#!/bin/sh
#
### BEGIN INIT INFO
# Provides: osaffmd
# Required-Start: $local_fs $remote_fs $network $syslog
# Required-Stop: $local_fs $remote_fs $network $syslog
# Should-Start: 
# Should-Stop: 
# Default-Start: 
# Default-Stop: 
# Short-Description: Start and stop the OpenSAF FM Server
# Description: OpenSAF FM Server Daemon
### END INIT INFO

osafdirfile=@sysconfdir@/@PACKAGE_NAME@/osafdir.conf
nodeinfofile=@sysconfdir@/plmcd.conf
prog="osaffmd"
initscript=`basename $0`

# Source LSB functions library
. /lib/lsb/init-functions

if [ ! -r $osafdirfile ]; then
	logger -t $prog "can't read $osafdirfile, exiting."
	exit 6
else
	. $osafdirfile
	. $pkgsysconfdir/nid.conf
	. $pkgsysconfdir/fmd.conf
fi	

if [ -f $nodeinfofile ]; then
        ee_id=`grep -i safEE $nodeinfofile`
        if [ $? -eq 0 ]; then
        export EE_ID=$ee_id
        fi
fi
         
binary=$pkglibdir/$prog
pidfile=$pkgpiddir/$prog.pid
lockfile=$lockdir/$initscript

RETVAL=0
NIDSERV="HLFM"
COMPNAMEFILE=$pkglocalstatedir/fmd_comp_name


start() {
	export LD_LIBRARY_PATH=$pkglibdir:$LD_LIBRARY_PATH
	[ -p $NIDFIFO ] || exit 1
	if [ ! -x $binary ]; then
		echo "$NID_MAGIC:$NIDSERV:$DAEMON_NOT_FND" > $NIDFIFO
		exit 5
	fi
	echo -n "Starting $prog: "
	start_daemon -p $pidfile $binary $args
	RETVAL=$?
	if [ $RETVAL -eq 0 ]; then
		touch $lockfile
		log_success_msg
	else
		log_failure_msg
		echo "$NID_MAGIC:$NIDSERV:$DAEMON_START_FAILED" > $NIDFIFO
	fi
	return $RETVAL
}

instantiate() {
	[ -p $COMPNAMEFILE ] || mkfifo -m 0644 $COMPNAMEFILE
	echo -n "AMF Instantiating $prog: "
	if [ ! -e $lockfile ]; then
		start_daemon -p $pidfile $binary $args
		RETVAL=$?
		sleep 1s
	fi
	echo -n "$SA_AMF_COMPONENT_NAME|$FMS_HA_ENV_HEALTHCHECK_KEY" > $COMPNAMEFILE &
        if [ $RETVAL -eq 0 ]; then
		log_success_msg
	else
		log_failure_msg
	fi
	return $RETVAL
}

stop() {
	echo -n "Stopping $prog: "
	killproc -p $pidfile $binary
	RETVAL=$?
	if [ $RETVAL -eq 0 ]; then
		rm -f $COMPNAMEFILE
		rm -f $lockfile
		log_success_msg
	else
		log_failure_msg
	fi
	return $RETVAL
}

restart() {
	stop
	start
}

tryrestart() {
	[ -e $lockfile ] && restart
	return 0
}

reload() {
	echo -n "Reloading $prog: "
	log_warning_msg
	logger -t $prog "reloading configuration is currently an unimplemented LSB feature"
	return 3
}

forcereload() {
	reload
	restart
	return 0
}

status() {
	pidofproc -p $pidfile $binary
	RETVAL=$?
	return $RETVAL
}

case "$1" in
  start)
	start
	RETVAL=$?
	;;
  instantiate)
	instantiate
	RETVAL=$?
	;;
  stop)
	stop
	RETVAL=$?
	;;
  restart)
	restart
	RETVAL=$?
	;;
  try-restart|condrestart)
	tryrestart
	RETVAL=$?
	;;
  reload)
	reload
	RETVAL=$?
	;;
  force-reload)
	forcereload
	RETVAL=$?
	;;
  status)
	status
	RETVAL=$?
	;;
  *)
	echo "Usage: $0 {start|instantiate|stop|restart|try-restart|reload|force-reload|status}"
	RETVAL=2
esac

exit $RETVAL
