# src/nid/opensafd.service.in
[Unit]
Description=OpenSAF daemon
After=remote-fs.target local-fs.target syslog.target network-online.target
Wants=network-online.target

[Service]
ExecStart=@initscriptdir@/opensafd start
ExecStop=@initscriptdir@/opensafd stop
Type=forking
PIDFile=@localstatedir@/run/opensaf/osafamfnd.pid
ControlGroup=cpu:/
TimeoutStartSec=3hours
KillMode=none
@systemdtasksmax@
Restart=on-failure

[Install]
WantedBy=multi-user.target
