# osaf/services/infrastructure/dtsv/dts/scripts/osaf-dtd.in
#!/bin/sh
#
### BEGIN INIT INFO
# Provides: osafdtd
# Required-Start: $local_fs $remote_fs $network $syslog
# Required-Stop: $local_fs $remote_fs $network $syslog
# Should-Start: 
# Should-Stop: 
# Default-Start: 
# Default-Stop: 
# Short-Description: Start and stop the OpenSAF DT Server
# Description: OpenSAF DT Server Daemon
### END INIT INFO

osafdirfile=@sysconfdir@/@PACKAGE_NAME@/osafdir.conf
prog="osafdtd"
initscript=`basename $0`

# Source LSB functions library
. /lib/lsb/init-functions

if [ ! -r $osafdirfile ]; then
	logger -t $prog "can't read $osafdirfile, exiting."
	exit 6
else
	. $osafdirfile
	. $pkgsysconfdir/nid.conf
	. $pkgsysconfdir/dtsv.conf
fi	

binary=$pkglibdir/$prog
pidfile=$pkgpiddir/$prog.pid
lockfile=$lockdir/$initscript

# Uncomment to enable tracing, or pas new command line options
#args="-T $prog.log"

RETVAL=0
NIDSERV="DTSV"
COMPNAMEFILE=$pkglocalstatedir/dts_comp_name
FILESFROM=/tmp/opensaf_dtsv

archives_rotate()
{
	archives_dir=$DTSV_ROOT_DIRECTORY/archives

	mkdir -p $archives_dir

	# Remove oldest archive
	oldest_arch=`ls $archives_dir/dtsv_*.tgz.$DTSV_MAX_ARCHIVES 2> /dev/null`
	if [ $? -eq 0 ]; then
		echo "  removing oldest archive `basename $oldest_arch`"
		rm -f $oldest_arch
	fi

	# Rename inbetween archives
	cnt=$DTSV_MAX_ARCHIVES
	while [ $cnt -gt 1 ]; do
		old=$cnt
		cnt=$(($cnt - 1))
		curr_file=`ls $archives_dir/dtsv_*.tgz.$cnt 2> /dev/null`
		if [ $? -eq 0 ]; then
			old_file=${curr_file%%tgz.$cnt}tgz.$old
			echo "  renaming `basename $curr_file` to `basename $old_file`"
			mv $curr_file $old_file
		fi
	done

	# Create new archive in background
	filename=$archives_dir/dtsv_`date +%G%m%d_%H%M%S`.tgz.1
	echo "  archiving current files into `basename $filename`"
	ionice -c3 tar -C $DTSV_ROOT_DIRECTORY -czf $filename --files-from=$FILESFROM --remove-files &
}

start() {
	export LD_LIBRARY_PATH=$pkglibdir:$LD_LIBRARY_PATH
	[ -p $NIDFIFO ] || exit 1
        if [ ! -x $binary ]; then
		echo "$NID_MAGIC:$NIDSERV:$DAEMON_NOT_FND" > $NIDFIFO
		exit 5
	fi

	export DTSV_COMP_NAME_FILE=$COMPNAMEFILE
	[ -d $DTSV_ROOT_DIRECTORY ] || mkdir -p $DTSV_ROOT_DIRECTORY

	#get the role of RDE to rotate the log files
	RDE_ROLE=`$bindir/rdegetrole`
	EXIT_STATUS=$?
	if [ $EXIT_STATUS -eq 0 ] ; then	
		if [ $RDE_ROLE = "ACTIVE" ] && [ $DTSV_MAX_ARCHIVES -gt 0 ]; then
			# Take a snapshot of the current set of files and store away for later use
			(cd $DTSV_ROOT_DIRECTORY; find -name "*.log" -amin +1) > $FILESFROM
		fi
	fi

	echo -n "Starting $prog: "
	start_daemon -p $pidfile $binary $args
	RETVAL=$?
	if [ $RETVAL -eq 0 ]; then
		touch $lockfile
		log_success_msg
	else
		log_failure_msg
		echo "$NID_MAGIC:$NIDSERV:$DAEMON_START_FAILED" > $NIDFIFO
	fi

	if [ $EXIT_STATUS -eq 0 ] ; then	
		if [ $RDE_ROLE = "ACTIVE" ] && [ $DTSV_MAX_ARCHIVES -gt 0 ]; then
			archives_rotate
		fi
	fi
	return $RETVAL
}

instantiate() {
	echo -n "AMF Instantiating $prog: "
	echo $SA_AMF_COMPONENT_NAME > $COMPNAMEFILE
	if [ -e $lockfile ]; then
		killproc -p $pidfile $binary -USR1
	else
		start_daemon -p $pidfile $binary $args
	fi
	RETVAL=$?
	if [ $RETVAL -eq 0 ]; then
		log_success_msg
	else
		log_failure_msg
	fi
	return $RETVAL
}

stop() {
	echo -n "Stopping $prog: "
	killproc -p $pidfile $binary
	RETVAL=$?
	if [ $RETVAL -eq 0 ]; then
		rm -f $COMPNAMEFILE
		rm -f $lockfile
		log_success_msg
	else
		log_failure_msg
	fi
	return $RETVAL
}

restart() {
	stop
	start
}

tryrestart() {
	[ -e $lockfile ] && restart
	return 0
}

reload() {
	echo -n "Reloading $prog: "
	log_warning_msg
	logger -t $prog "reloading configuration is currently an unimplemented LSB feature"
	return 3
}

forcereload() {
	reload
	restart
	return 0
}

status() {
	pidofproc -p $pidfile $binary
	RETVAL=$?
	return $RETVAL
}

case "$1" in
  start)
	start
	RETVAL=$?
	;;
  instantiate)
	instantiate
	RETVAL=$?
	;;
  stop)
	stop
	RETVAL=$?
	;;
  restart)
	restart
	RETVAL=$?
	;;
  try-restart|condrestart)
	tryrestart
	RETVAL=$?
	;;
  reload)
	reload
	RETVAL=$?
	;;
  force-reload)
	forcereload
	RETVAL=$?
	;;
  status)
	status
	RETVAL=$?
	;;
  *)
	echo "Usage: $0 {start|instantiate|stop|restart|try-restart|reload|force-reload|status}"
	RETVAL=2
esac

exit $RETVAL
